# config/wsdd.service
[Unit]
Description=Web Services Dynamic Discovery host daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/wsdd --workgroup WORKGROUP --hostname BANANAPI --chroot /run/wsdd
Restart=always
RestartSec=10
User=nobody
Group=nogroup
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=wsdd
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target